# tools/ivr_beacon/launchd/com.ivr.probe.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Generated by setup.sh — HOME_DIR is substituted at install time -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ivr.probe.hourly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>HOME_DIR/ivr/run_ivr_batch.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TWILIO_ACCOUNT_SID</key>
        <string>TWILIO_SID_VALUE</string>
        <key>TWILIO_AUTH_TOKEN</key>
        <string>TWILIO_TOKEN_VALUE</string>
        <key>TWILIO_PHONE_NUMBER</key>
        <string>TWILIO_NUMBER_VALUE</string>
    </dict>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/ivr_probe.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ivr_probe.log</string>
</dict>
</plist>
